# conf/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /var/django/aws-eb-docker-django-skeleton/conf/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
